# rootfs/etc/s6/radicale/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

STARTCMD=""
exec ${STARTCMD}
