# rootfs/etc/s6/clamav/run
#!/bin/bash

if [ -f ./setup ]
then
  ./setup
fi

STARTCMD="/usr/sbin/amavisd foreground"
exec ${STARTCMD}
