# rootfs/etc/s6/amavisd/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

STARTCMD="/usr/sbin/amavisd foreground"
exec ${STARTCMD}
